# rootfs/etc/supervisor.d/dashing.ini

[program:dashing]
directory      = /opt/icinga2
command        = /usr/bin/thin --pid /tmp/thin.pid --log /tmp/thin.log --port 3030 --environment production --rackup /opt/icinga2/config.ru start
stdout_logfile = /tmp/%(program_name)s.log
stderr_logfile = /tmp/%(program_name)s.log
autorestart    = true
